# roles/backup/files/personal-backup-maintenance.timer
[Unit]
Description=Run maintenance tasks for personal backup once a day

[Timer]
OnCalendar=*-*-* 23:45:00

[Install]
WantedBy=timers.target
